# debian/postinst
#!/bin/sh
# postinst script for mangezmieux
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

# Use debconf.
. /usr/share/debconf/confmodule

install_dir="/usr/share/mangezmieux"
dbinstall="${install_dir}/scripts/db_install.sh"
settings="${install_dir}/settings.py"

case "$1" in
    configure)
        # Install python deps not in debian repos
        pip install django-filter djangorestframework

        # Use debconf to ask DB credentials
        db_get mangezmieux/mangezmieux_db_host && db_host="$RET"
        db_get mangezmieux/mangezmieux_db_name && db_name="$RET"
        db_get mangezmieux/mangezmieux_db_login && db_login="$RET"
        db_get mangezmieux/mangezmieux_db_passwd && db_passwd="$RET"
        db_get mangezmieux/mangezmieux_install_db && install_req="$RET"

        # Changing parameters in settings.py with debconf results
        sed -i "s/'HOST':\ '.*'/'HOST':\ '${db_host}'/" ${settings}
        sed -i "s/'NAME':\ '.*'/'NAME':\ '${db_name}'/" ${settings}
        sed -i "s/'USER':\ '.*'/'USER':\ '${db_login}'/" ${settings}
        sed -i "s/'PASSWORD':\ '.*'/'PASSWORD':\ '${db_passwd}'/" ${settings}
        sed -i 's/DEBUG\ =\ True/DEBUG\ = \False/' ${settings}

        # Install DB
        if [ "${install_req}" = "true" ]; then
            ${dbinstall} -H ${db_host} -d ${db_name} -u ${db_login} -p ${db_passwd}
        fi
        
        #Upload dir is rwx for everyone
        chown -R www-data ${install_dir}
        chmod a+rwx ${install_dir}/upload

        # Configure Apache with WSGI
        a2ensite mangezmieux
        invoke-rc.d apache2 reload
        exit 0
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#D
